# systemd/prune-logs.service
[Unit]
Description=Prune MESA logs

[Service]
Type=simple
ExecStart=/var/www/mesa-logs/systemd/prune-logs.sh

[Install]
WantedBy=default.target
